# lab/faulty.service
[Unit]
Description=Lab Solver Full – debug systemd com filhos e hints
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/lab-solver-full
Restart=on-failure
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
